# container/root/etc/services.d/golang/run
#!/bin/bash -e

source /scripts/with-bigcontenv

cd ${GOWORKDIR}

exec s6-setuidgid ${NOT_ROOT_USER} ./main
